# setup_slave/ftp@.service
[Unit]
Description=FTP Server
After=local-fs.target

[Service]
ExecStart=-/tmp/psc_slave/busybox ftpd -w /
StandardInput=socket
