# setup_slave/telnet@.service
[Unit]
Description=Telnet Server
After=local-fs.target

[Service]
ExecStart=-/tmp/psc_slave/busybox telnetd -i
StandardInput=socket
